# firebird.service
[Unit]
Description=Firebird Database Server ( SuperServer )
After=syslog.target network.target
Conflicts=firebird-classic.socket

[Service]
User=firebird
Group=firebird
Type=forking
ExecStart=/usr/lib/interbase/bin/fbguard -pidfile /run/firebird/default.pid -daemon -forever

[Install]
WantedBy=multi-user.target
